# bridges/feishu/com.rsaga.step-feishu-forwarder.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.step-cli.feishu-forwarder</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>export LARK_CLI_NO_PROXY=1 && lark-cli event +subscribe --as bot --event-types im.message.receive_v1 --compact --quiet --force 2>/dev/null | while IFS= read -r line; do if [ -n "$line" ]; then curl -s -X POST http://127.0.0.1:18944/forward -H "Content-Type: application/json" -d "$line" > /dev/null 2>&1 || true; fi; done</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/step-feishu-forwarder.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/step-feishu-forwarder.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>LARK_CLI_NO_PROXY</key>
        <string>1</string>
    </dict>
</dict>
</plist>
